# instance-builder/aws/ec2/api/api.service
[Unit]
Description=VLLM Server
After=network.target

[Service]
ExecStart=/home/ubuntu/api --config /home/ubuntu/config.toml
Restart=always
User=ubuntu
Group=ubuntu

[Install]
WantedBy=multi-user.target

